# assets/etc/service/telnetd/run
#!/bin/sh
cd /
exec chpst -P /usr/sbin/telnetd -F -l /bin/sh
