# assets/etc/service/tty1/run
#!/bin/sh
cd /
exec chpst -P /sbin/getty -n -l /bin/ash 0 tty1
